# .config/systemd/user/nvchecker.timer
[Unit]
Description=Daily check for software updates

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
